# kurage-systemd/kfreqai-kcbrain-shadow.service
[Unit]
Description=Kurage Freq AI Trade - hourly kcbrain shadow regime/ranking (no live effect)

[Service]
Type=oneshot
WorkingDirectory=/home/kojima/work/kfreqai/kurage-advisory
ExecStart=/usr/bin/python3 kcbrain_shadow.py
TimeoutStartSec=600
