# kurage-systemd/kfreqai-regime-hourly.timer
[Unit]
Description=Run kfreqai hourly gemma4 regime tagger

[Timer]
OnCalendar=hourly
RandomizedDelaySec=120
Persistent=true

[Install]
WantedBy=timers.target
