# kurage-systemd/kfreqai-directive-daily.service
[Unit]
Description=Kurage Freq AI Trade - Claude CLI risk directive + blog post (3x/day)

[Service]
Type=oneshot
WorkingDirectory=/home/kojima/work/kfreqai/kurage-advisory
EnvironmentFile=/home/kojima/work/aixec/.env
ExecStart=/usr/bin/python3 daily_directive.py
ExecStart=/usr/bin/python3 blog_post.py
TimeoutStartSec=600
